# deploy/com.kalshi.research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Local daily research loop (uses your Max-plan claude login). Install:
       cp deploy/com.kalshi.research.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.kalshi.research.plist
     Runs ~10:00 local daily; logs to /tmp/kalshi_research.{out,err}. If the Mac is
     asleep at the scheduled time, launchd runs it at next wake. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.kalshi.research</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sayujjain/Documents/Workspace/Prodiction Market Trends/deploy/research_local.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/tmp/kalshi_research.out</string>
  <key>StandardErrorPath</key><string>/tmp/kalshi_research.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
